# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ============================================
; Laravel Queues (if using)
; ============================================
[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work --queue=default,high --timeout=900 --tries=3 --max-jobs=1000
autostart=true
autorestart=true
numprocs=4
redirect_stderr=true
stdout_logfile=/var/log/laravel-worker.log
stopwaitsecs=3600

; ============================================
; PHP-FPM
; ============================================
[program:php-fpm]
command=/usr/sbin/php-fpm8.2 -y /usr/local/etc/php-fpm.conf -R
process_name=php-fpm
numprocs=1
priority=999
autostart=true
autorestart=true
startsecs=0
redirect_stderr=true
stdout_logfile=/var/log/php-fpm.log

; ============================================
; Nginx
; ============================================
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
process_name=nginx
numprocs=1
priority=999
autostart=true
autorestart=true
startsecs=0
redirect_stderr=true
stdout_logfile=/var/log/nginx.log

; ============================================
; Horizon (if using Laravel Horizon)
; ============================================
; [program:horizon]
; process_name=horizon
; command=php /var/www/artisan horizon
; autostart=true
; autorestart=true
; numprocs=1
; redirect_stderr=true
; stdout_logfile=/var/log/horizon.log
; stopwaitsecs=3600

; ============================================
; Cron Jobs (if needed)
; ============================================
; [program:schedule]
; command=/bin/bash -c "while true; do php /var/www/artisan schedule:run >> /var/log/cron.log 2>&1; sleep 60; done"
; autostart=true
; autorestart=true
; numprocs=1
; redirect_stderr=true
; stdout_logfile=/var/log/cron.log
